# hik-server.service
[Unit]
Description=Hikvision Server Service Files

[Service]
ExecStart=python3 /home/james/hik/server.py

[Install]
WantedBy=multi-user.target
